# scripts/k3s/beat-monitor.service
# Systemd service file for Beat health monitoring
# Copy to: /etc/systemd/system/beat-monitor.service
# Enable: systemctl enable beat-monitor && systemctl start beat-monitor

[Unit]
Description=Celery Beat Health Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app
ExecStart=/app/skyeye/scripts/k3s/monitor_beat_health.sh
Restart=always
RestartSec=10
StandardOutput=append:/app/logs/beat-monitor.log
StandardError=append:/app/logs/beat-monitor.log

# 防止监控脚本自身异常
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target